# systemd/gpioctld.socket
[Unit]
Description=gpioctld UNIX socket

[Socket]
ListenSequentialPacket=/run/gpioctld.socket

[Install]
WantedBy=sockets.target
